# containers/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
/usr/sbin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile